# supervisord.conf
[program:drission-rod]
command=uv run /workspace/main.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
environment=DISPLAY=":99"
priority=99